# data/global/script/cron/0unison.cron
#!/bin/bash

# Action!
i="2"
# do two times
while [ $i -gt 0 ]
do
	echo "UNISON RUN $i"
	sudo unison default.prf
	sudo unison log.prf
	sudo unison unison.prf
	i=$(( $i - 1 ))
	sleep 10 || exit 2 
done

exit 0